# ext-sources/consolekit.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
    Copyright 2008 Sun Microsystems, Inc.  All rights reserved.
    Use is subject to license terms.

    pragma ident	"@(#)server.xml	1.2	04/08/09 SMI"
-->

<service_bundle type='manifest' name='SFEconsolekit:consoleKit'>

<service
    name='system/consolekit'
    type='service'
    version='1'>

	<dependency
	    name='filesystem'
	    grouping='require_all'
	    restart_on='none'
	    type='service'>
		<service_fmri value='svc:/system/filesystem/local'/>
	</dependency>

        <dependency
            name='dbus'
            type='service'
            grouping='require_all'
            restart_on='none'>
                <service_fmri value='svc:/system/dbus' />
        </dependency>

	<exec_method
	    type='method'
	    name='start'
	    exec='/usr/sbin/console-kit-daemon'
	    timeout_seconds='60' />

	<exec_method
	    type='method'
	    name='stop'
	    exec=':kill'
	    timeout_seconds='60' />

	<instance name='default' enabled='false' />

	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
				ConsoleKit
			</loctext>
		</common_name>
		<documentation>
			<manpage title='console-kit-daemon' section='3'
				manpath='/usr/share/man' />
		</documentation>
	</template>
</service>

</service_bundle>
